# debian/update-manager-core.cron.daily
#!/bin/sh -e

/usr/lib/update-manager/check-new-release > /var/lib/update-notifier/release-upgrade-available
